# config/systemd/vantus.service
[Unit]
Description=Vantus Next.js Production Application (Fortune-500 Hardened)
Documentation=https://vantus.systems/docs
After=network-online.target nginx.service
Wants=network-online.target
Requires=nginx.service

[Service]
Type=simple
User=vantus
Group=vantus
WorkingDirectory=/var/www/vantus
EnvironmentFile=-/var/www/vantus/.env

# State, Log, and Runtime Directories with secure permissions
StateDirectory=vantus
StateDirectoryMode=750
LogsDirectory=vantus
LogsDirectoryMode=750
RuntimeDirectory=vantus
RuntimeDirectoryMode=750

# Environment with secure defaults
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3005
Environment=UPLOADS_DIR=/var/lib/vantus/uploads
Environment=TZ=UTC

# Comprehensive Security Hardening
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
MemoryDenyWriteExecute=true
LockPersonality=true

# Resource Limits for production stability
LimitNOFILE=65536
LimitNPROC=1024
LimitCORE=0
LimitMEMLOCK=64M
LimitAS=infinity
LimitRSS=infinity
LimitFSIZE=infinity

# Capability Bounding Set (minimal required)
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Secure execution
ExecStartPre=/usr/bin/mkdir -p /var/lib/vantus/uploads
ExecStartPre=/usr/bin/chown vantus:vantus /var/lib/vantus/uploads
ExecStartPre=/usr/bin/chmod 750 /var/lib/vantus/uploads

# Main application start
ExecStart=/usr/bin/npm start

# Restart policy with exponential backoff
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=300
StartLimitBurst=5
StartLimitAction=reboot-force

# Process management
TimeoutSec=300
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vantus-app
SyslogLevelPrefix=true

# Performance optimization
Nice=0
IOSchedulingClass=best-effort
IOSchedulingPriority=4
CPUSchedulingPolicy=other
CPUSchedulingPriority=0

# Health monitoring
WatchdogSec=60
RuntimeMaxSec=86400

[Install]
WantedBy=multi-user.target
