# config/systemd/vantus.service
[Unit]
Description=crispy-fishstick Next.js Application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vantus
Group=vantus
WorkingDirectory=/var/www/vantus
EnvironmentFile=/var/www/vantus/.env

# Create state, logs, and runtime directories
StateDirectory=vantus
LogsDirectory=vantus
RuntimeDirectory=vantus

# Environment
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3005

# Security Hardening
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
RestrictRealtime=true
RestrictNamespaces=true

# Pre-start: Run migrations
# Note: We use npx to run the local prisma binary
# ExecStartPre=/usr/bin/npx prisma migrate deploy
# ExecStartPre=/usr/bin/npm run db:migrate

# Start the application
# Assumes 'npm run build' was run during deployment
ExecStart=/usr/bin/npm start

# Restart policy
Restart=on-failure
RestartSec=10
# StartLimitInterval=200
# StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vantus

# Performance
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
